# mysql/visord_mysql.conf
[supervisord]
nodaemon=true
[program:sshd]
command = sh -c "rc-service mariadb start && kill -s SIGTERM $(cat supervisord.pid)"
